# rserver_git.timer
[Unit]
Description=RSERVER_GIT

[Timer]
OnStartupSec=1min
OnCalendar=1h
