# rserver_git.timer
[Unit]
Description=RSERVER_GIT
Requires=rserver_git.service

[Timer]
OnStartupSec=1min
OnCalendar=3h
AccuracySec=1min
Unit=rserver_git.service
