# rserver_git.service
[Unit]
Description=RSERVER_GIT

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=
ExecStart=

#Restart=always
#RestartSec=3
#TimeoutSec=30
